# config/systemd/user/hypridle.service.d/override.conf
[Unit]
Conflicts=hypridle-awake.service
