# config/systemd/user/dunst.service.d/override.conf
[Install]
WantedBy=graphical-session.target
