# etc/systemd/user/eww-daemon.service
[Unit]
Description=ElKowars wacky widgets daemon
Documentation=https://elkowar.github.io/eww/
PartOf=graphical-session.target

[Service]
Type=exec
ExecStart=/usr/bin/eww daemon --no-daemonize
Restart=on-failure

[Install]
WantedBy=graphical-session.target
