# etc/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland session 
Requisite=default.target
After=default.target
