# config/systemd/user/foot-server@.socket.d/override.conf
[Install]
WantedBy=hyprland-session.target
